# mesh/etc/rc.local

/usr/bin/sysupgrade.sh > /tmp/log.txt &
